# dist/amazon-linux/2/after-remove.sh
#!/bin/bash

set -ex

if systemctl is-enabled --quiet metrist-orchestrator; then
  systemctl stop metrist-orchestrator
fi

systemctl daemon-reload

