# dist/ubuntu/20.04/inc/etc/systemd/system/metrist-agent.service
[Unit]
Description=Metrist.io Agent

[Service]
ExecStart=/usr/bin/metrist-agent
EnvironmentFile=-/etc/default/metrist-agent

[Install]
WantedBy=multi-user.target
