# dist/ubuntu/20.04/inc/etc/systemd/system/metrist-orchestrator.service
[Unit]
Description=Metrist.io Orchestrator service

[Service]
ExecStart=/usr/bin/metrist-orchestrator
EnvironmentFile=-/etc/default/metrist-orchestrator

[Install]
WantedBy=multi-user.target
